# job-catalog-automation/com.user.jobcatalog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.{{USER_SYSTEM_USERNAME}}.jobcatalog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/{{USER_SYSTEM_USERNAME}}/Documents/Projects/Personal Agent/job-catalog-automation/job_catalog_updater.js</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/com.{{USER_SYSTEM_USERNAME}}.jobcatalog.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.{{USER_SYSTEM_USERNAME}}.jobcatalog.err</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
